# rootfs/usr/lib/systemd/system/apport-autoreport.path
[Unit]
Description=Process error reports when automatic reporting is enabled (file watch)
ConditionPathExists=/var/lib/apport/autoreport

[Path]
PathChanged=/var/crash

[Install]
WantedBy=paths.target
