# rootfs/var/lib/dpkg/info/apport.conffiles
/etc/apport/crashdb.conf
/etc/apport/report-ignore/README.denylist
/etc/cron.daily/apport
/etc/default/apport
/etc/logrotate.d/apport
