# rootfs/var/lib/dpkg/info/dbus.list
/.
/etc
/etc/default
/etc/default/dbus
/etc/init.d
/etc/init.d/dbus
/usr
/usr/lib
/usr/lib/dbus-1.0
/usr/lib/dbus-1.0/dbus-daemon-launch-helper
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/dbus.service
/usr/lib/systemd/system/multi-user.target.wants
/usr/share
/usr/share/bug
/usr/share/bug/dbus
/usr/share/bug/dbus/control
/usr/share/doc
/usr/share/doc/dbus
/usr/share/doc/dbus/README.Debian
/usr/share/doc/dbus/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/dbus
/var
/var/lib
/var/lib/dbus
/usr/lib/systemd/system/multi-user.target.wants/dbus.service
/usr/share/doc/dbus/AUTHORS.gz
/usr/share/doc/dbus/NEWS.gz
/usr/share/doc/dbus/README.gz
/usr/share/doc/dbus/changelog.Debian.gz
