# rootfs/var/lib/dpkg/info/dbus-system-bus-common.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/dbus.socket
/usr/lib/systemd/system/sockets.target.wants
/usr/lib/sysusers.d
/usr/lib/sysusers.d/dbus.conf
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system.conf
/usr/share/doc
/usr/share/doc/dbus-system-bus-common
/usr/share/doc/dbus-system-bus-common/AUTHORS.gz
/usr/share/doc/dbus-system-bus-common/NEWS.gz
/usr/share/doc/dbus-system-bus-common/README.gz
/usr/share/doc/dbus-system-bus-common/changelog.Debian.gz
/usr/share/doc/dbus-system-bus-common/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/dbus-system-bus-common
/usr/lib/systemd/system/sockets.target.wants/dbus.socket
