# rootfs/var/lib/dpkg/info/kmscon.list
/.
/etc
/etc/needrestart
/etc/needrestart/conf.d
/etc/needrestart/conf.d/kmscon.conf
/usr
/usr/bin
/usr/bin/kmscon
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/kmscon
/usr/lib/arm-linux-gnueabihf/kmscon/mod-bbulk.so
/usr/lib/arm-linux-gnueabihf/kmscon/mod-gltex.so
/usr/lib/arm-linux-gnueabihf/kmscon/mod-pango.so
/usr/lib/arm-linux-gnueabihf/kmscon/mod-pixman.so
/usr/lib/arm-linux-gnueabihf/kmscon/mod-unifont.so
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/kmsconvt@.service
/usr/libexec
/usr/libexec/kmscon
/usr/libexec/kmscon/kmscon
/usr/share
/usr/share/doc
/usr/share/doc/kmscon
/usr/share/doc/kmscon/changelog.Debian.gz
/usr/share/doc/kmscon/copyright
/usr/share/doc/kmscon/unicode-test.txt
/usr/share/doc/kmscon/vte.txt
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/kmscon
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/kmscon.1.gz
