# rootfs/var/lib/dpkg/info/libgl1-amber-dri:armhf.list
/.
/usr
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/dri
/usr/lib/arm-linux-gnueabihf/dri/nouveau_vieux_dri.so
/usr/lib/arm-linux-gnueabihf/dri/r200_dri.so
/usr/lib/arm-linux-gnueabihf/dri/radeon_dri.so
/usr/share
/usr/share/bug
/usr/share/bug/libgl1-amber-dri
/usr/share/bug/libgl1-amber-dri/control
/usr/share/bug/libgl1-amber-dri/script
/usr/share/doc
/usr/share/doc/libgl1-amber-dri
/usr/share/doc/libgl1-amber-dri/changelog.Debian.gz
/usr/share/doc/libgl1-amber-dri/copyright
/usr/share/drirc.d
/usr/share/drirc.d/00-amber-defaults.conf
