# rootfs/var/lib/dpkg/info/libnss-ldap:armhf.list
/.
/etc
/etc/init.d
/etc/init.d/libnss-ldap
/lib
/lib/arm-linux-gnueabihf
/lib/arm-linux-gnueabihf/libnss_ldap-.so
/usr
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/sbin
/usr/sbin/nssldap-update-ignoreusers
/usr/share
/usr/share/doc
/usr/share/doc/libnss-ldap
/usr/share/doc/libnss-ldap/AUTHORS
/usr/share/doc/libnss-ldap/LDAP-Permissions.txt.gz
/usr/share/doc/libnss-ldap/NEWS.gz
/usr/share/doc/libnss-ldap/README.Debian
/usr/share/doc/libnss-ldap/README.gz
/usr/share/doc/libnss-ldap/changelog.Debian.gz
/usr/share/doc/libnss-ldap/copyright
/usr/share/doc/libnss-ldap/examples
/usr/share/doc/libnss-ldap/examples/groups.ldif
/usr/share/doc/libnss-ldap/examples/ldap.conf.gz
/usr/share/doc/libnss-ldap/examples/nsswitch.ldap
/usr/share/doc/libnss-ldap/examples/people.ldif
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/libnss-ldap
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/nss_ldap.5.gz
/usr/share/man/man8
/usr/share/man/man8/nssldap-update-ignoreusers.8.gz
/var
/var/lib
/var/lib/libnss-ldap
/lib/arm-linux-gnueabihf/libnss_ldap.so.2
/usr/lib/arm-linux-gnueabihf/libnss_ldap.so
