# rootfs/var/lib/dpkg/info/mesa-vulkan-drivers:armhf.list
/.
/usr
/usr/bin
/usr/bin/mesa-overlay-control.py
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/libVkLayer_MESA_device_select.so
/usr/lib/arm-linux-gnueabihf/libVkLayer_MESA_overlay.so
/usr/lib/arm-linux-gnueabihf/libvulkan_broadcom.so
/usr/lib/arm-linux-gnueabihf/libvulkan_freedreno.so
/usr/lib/arm-linux-gnueabihf/libvulkan_lvp.so
/usr/lib/arm-linux-gnueabihf/libvulkan_radeon.so
/usr/lib/arm-linux-gnueabihf/libvulkan_virtio.so
/usr/share
/usr/share/bug
/usr/share/bug/mesa-vulkan-drivers
/usr/share/bug/mesa-vulkan-drivers/control
/usr/share/bug/mesa-vulkan-drivers/script
/usr/share/doc
/usr/share/doc/mesa-vulkan-drivers
/usr/share/doc/mesa-vulkan-drivers/changelog.Debian.gz
/usr/share/doc/mesa-vulkan-drivers/copyright
/usr/share/drirc.d
/usr/share/drirc.d/00-radv-defaults.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mesa-vulkan-drivers
/usr/share/vulkan
/usr/share/vulkan/explicit_layer.d
/usr/share/vulkan/explicit_layer.d/VkLayer_MESA_overlay.json
/usr/share/vulkan/icd.d
/usr/share/vulkan/icd.d/broadcom_icd.armv7l.json
/usr/share/vulkan/icd.d/freedreno_icd.armv7l.json
/usr/share/vulkan/icd.d/lvp_icd.armv7l.json
/usr/share/vulkan/icd.d/radeon_icd.armv7l.json
/usr/share/vulkan/icd.d/virtio_icd.armv7l.json
/usr/share/vulkan/implicit_layer.d
/usr/share/vulkan/implicit_layer.d/VkLayer_MESA_device_select.json
