# rootfs/var/lib/dpkg/info/nscd.list
/.
/etc
/etc/init.d
/etc/init.d/nscd
/etc/nscd.conf
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/nscd.service
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/nscd.conf
/usr/sbin
/usr/sbin/nscd
/usr/share
/usr/share/doc
/usr/share/doc/nscd
/usr/share/doc/nscd/changelog.Debian.gz
/usr/share/doc/nscd/changelog.gz
/usr/share/doc/nscd/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/nscd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/nscd.8.gz
/var
/var/cache
/var/cache/nscd
