# rootfs/var/lib/dpkg/info/nslcd.conffiles
/etc/default/nslcd
/etc/init.d/nslcd
/etc/network/if-up.d/nslcd
