# rootfs/var/lib/dpkg/info/nslcd.list
/.
/etc
/etc/default
/etc/default/nslcd
/etc/init.d
/etc/init.d/nslcd
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/nslcd
/usr
/usr/lib
/usr/lib/networkd-dispatcher
/usr/lib/networkd-dispatcher/routable.d
/usr/sbin
/usr/sbin/nslcd
/usr/share
/usr/share/doc
/usr/share/doc/nslcd
/usr/share/doc/nslcd/AUTHORS.gz
/usr/share/doc/nslcd/NEWS.gz
/usr/share/doc/nslcd/README.gz
/usr/share/doc/nslcd/changelog.Debian.gz
/usr/share/doc/nslcd/copyright
/usr/share/doc/nslcd/examples
/usr/share/doc/nslcd/examples/ldapns.ldif
/usr/share/doc/nslcd/examples/ldapns.schema
/usr/share/doc/nslcd/examples/nslcd.conf
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/nslcd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/nslcd.8.gz
/usr/lib/networkd-dispatcher/routable.d/nslcd
