# rootfs/var/lib/dpkg/info/seatd.list
/.
/etc
/etc/default
/etc/default/seatd
/etc/init.d
/etc/init.d/seatd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/seatd.service
/usr
/usr/bin
/usr/bin/seatd-launch
/usr/sbin
/usr/sbin/seatd
/usr/share
/usr/share/doc
/usr/share/doc/seatd
/usr/share/doc/seatd/README.md
/usr/share/doc/seatd/TODO.Debian
/usr/share/doc/seatd/changelog.Debian.gz
/usr/share/doc/seatd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/seatd-launch.1.gz
/usr/share/man/man8
/usr/share/man/man8/seatd.8.gz
