# rootfs/var/lib/dpkg/info/sudo.list
/.
/etc
/etc/pam.d
/etc/pam.d/sudo
/etc/pam.d/sudo-i
/etc/sudo.conf
/etc/sudo_logsrvd.conf
/etc/sudoers
/etc/sudoers.d
/etc/sudoers.d/README
/usr
/usr/bin
/usr/bin/cvtsudoers
/usr/bin/sudo
/usr/bin/sudoreplay
/usr/include
/usr/include/sudo_plugin.h
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/sudo.conf
/usr/libexec
/usr/libexec/sudo
/usr/libexec/sudo/audit_json.so
/usr/libexec/sudo/group_file.so
/usr/libexec/sudo/libsudo_util.so.0.0.0
/usr/libexec/sudo/sesh
/usr/libexec/sudo/sudo_intercept.so
/usr/libexec/sudo/sudo_noexec.so
/usr/libexec/sudo/sudoers.so
/usr/libexec/sudo/system_group.so
/usr/sbin
/usr/sbin/sudo_logsrvd
/usr/sbin/sudo_sendlog
/usr/sbin/visudo
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_sudo.py
/usr/share/doc
/usr/share/doc/sudo
/usr/share/doc/sudo/CONTRIBUTING.md
/usr/share/doc/sudo/CONTRIBUTORS.md.gz
/usr/share/doc/sudo/HISTORY.md
/usr/share/doc/sudo/NEWS.Debian.gz
/usr/share/doc/sudo/NEWS.gz
/usr/share/doc/sudo/OPTIONS
/usr/share/doc/sudo/README.Debian
/usr/share/doc/sudo/README.md
/usr/share/doc/sudo/SECURITY.md
/usr/share/doc/sudo/TROUBLESHOOTING.md.gz
/usr/share/doc/sudo/UPGRADE.md.gz
/usr/share/doc/sudo/changelog.Debian.gz
/usr/share/doc/sudo/copyright
/usr/share/doc/sudo/examples
/usr/share/doc/sudo/examples/cvtsudoers.conf
/usr/share/doc/sudo/examples/pam.conf
/usr/share/doc/sudo/examples/sudo.conf
/usr/share/doc/sudo/examples/sudo_logsrvd.conf
/usr/share/doc/sudo/examples/sudoers
/usr/share/doc/sudo/examples/sudoers.dist
/usr/share/doc/sudo/examples/syslog.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/sudo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cvtsudoers.1.gz
/usr/share/man/man5
/usr/share/man/man5/sudo.conf.5.gz
/usr/share/man/man5/sudo_logsrv.proto.5.gz
/usr/share/man/man5/sudo_logsrvd.conf.5.gz
/usr/share/man/man5/sudo_plugin.5.gz
/usr/share/man/man5/sudoers.5.gz
/usr/share/man/man5/sudoers_timestamp.5.gz
/usr/share/man/man8
/usr/share/man/man8/sudo.8.gz
/usr/share/man/man8/sudo_logsrvd.8.gz
/usr/share/man/man8/sudo_root.8.gz
/usr/share/man/man8/sudo_sendlog.8.gz
/usr/share/man/man8/sudoreplay.8.gz
/usr/share/man/man8/visudo.8.gz
/var
/var/lib
/var/lib/sudo
/var/lib/sudo/lectured
/usr/bin/sudoedit
/usr/lib/systemd/system/sudo.service
/usr/libexec/sudo/libsudo_util.so
/usr/libexec/sudo/libsudo_util.so.0
/usr/share/man/man8/sudoedit.8.gz
