# rootfs/var/lib/dpkg/info/systemd-resolved.list
/.
/etc
/etc/systemd
/etc/systemd/resolved.conf
/usr
/usr/bin
/usr/bin/resolvectl
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/systemd-resolved.service
/usr/lib/systemd/systemd-resolved
/usr/lib/sysusers.d
/usr/lib/sysusers.d/systemd-resolve.conf
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/systemd-resolve.conf
/usr/sbin
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/resolvectl
/usr/share/bash-completion/completions/systemd-resolve
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.resolve1.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.resolve1.conf
/usr/share/doc
/usr/share/doc/systemd-resolved
/usr/share/doc/systemd-resolved/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/systemd-resolved
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/resolvectl.1.gz
/usr/share/man/man5
/usr/share/man/man5/org.freedesktop.resolve1.5.gz
/usr/share/man/man5/resolved.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/systemd-resolved.service.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.resolve1.policy
/usr/share/zsh
/usr/share/zsh/vendor-completions
/usr/share/zsh/vendor-completions/_resolvectl
/usr/sbin/resolvconf
/usr/share/doc/systemd-resolved/NEWS.Debian.gz
/usr/share/doc/systemd-resolved/changelog.Debian.gz
/usr/share/man/man1/resolvconf.1.gz
/usr/share/man/man5/resolved.conf.d.5.gz
/usr/share/man/man8/systemd-resolved.8.gz
