# rootfs/var/lib/dpkg/info/ubuntu-pro-client.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/ubuntu_pro_apt_news
/etc/apparmor.d/ubuntu_pro_esm_cache
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/20apt-esm-hook.conf
/etc/apt/preferences.d
/etc/apt/preferences.d/ubuntu-pro-esm-apps
/etc/apt/preferences.d/ubuntu-pro-esm-infra
/etc/logrotate.d
/etc/logrotate.d/ubuntu-pro-client
/etc/ubuntu-advantage
/etc/ubuntu-advantage/uaclient.conf
/etc/update-manager
/etc/update-manager/release-upgrades.d
/etc/update-manager/release-upgrades.d/ubuntu-advantage-upgrades.cfg
/etc/update-motd.d
/etc/update-motd.d/91-contract-ua-esm-status
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/apt-news.service
/lib/systemd/system/esm-cache.service
/lib/systemd/system/ua-reboot-cmds.service
/lib/systemd/system/ua-timer.service
/lib/systemd/system/ua-timer.timer
/lib/systemd/system/ubuntu-advantage.service
/usr
/usr/bin
/usr/bin/ubuntu-advantage
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/uaclient
/usr/lib/python3/dist-packages/uaclient/__init__.py
/usr/lib/python3/dist-packages/uaclient/actions.py
/usr/lib/python3/dist-packages/uaclient/api
/usr/lib/python3/dist-packages/uaclient/api/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/api.py
/usr/lib/python3/dist-packages/uaclient/api/data_types.py
/usr/lib/python3/dist-packages/uaclient/api/errors.py
/usr/lib/python3/dist-packages/uaclient/api/exceptions.py
/usr/lib/python3/dist-packages/uaclient/api/u
/usr/lib/python3/dist-packages/uaclient/api/u/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/apt_news
/usr/lib/python3/dist-packages/uaclient/api/u/apt_news/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/apt_news/current_news
/usr/lib/python3/dist-packages/uaclient/api/u/apt_news/current_news/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/apt_news/current_news/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro
/usr/lib/python3/dist-packages/uaclient/api/u/pro/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/configure_retry_service
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/configure_retry_service/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/configure_retry_service/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/full_auto_attach
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/full_auto_attach/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/full_auto_attach/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/should_auto_attach
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/should_auto_attach/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/auto/should_auto_attach/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/initiate
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/initiate/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/initiate/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/revoke
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/revoke/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/revoke/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/wait
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/wait/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/magic/wait/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/token
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/token/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/token/full_token_attach
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/token/full_token_attach/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/attach/token/full_token_attach/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/detach
/usr/lib/python3/dist-packages/uaclient/api/u/pro/detach/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/detach/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/summary
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/summary/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/summary/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/updates
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/updates/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/packages/updates/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/execute
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/execute/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/execute/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/plan
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/plan/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/_common/plan/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/execute
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/execute/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/execute/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/plan
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/plan/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/cve/plan/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/execute
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/execute/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/execute/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/plan
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/plan/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/fix/usn/plan/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/livepatch_cves
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/livepatch_cves/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/livepatch_cves/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/reboot_required
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/reboot_required/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/security/status/reboot_required/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/dependencies
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/dependencies/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/dependencies/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/disable
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/disable/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/disable/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/enable
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/enable/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/services/enable/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/enabled_services
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/enabled_services/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/enabled_services/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/is_attached
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/is_attached/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/status/is_attached/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/version
/usr/lib/python3/dist-packages/uaclient/api/u/pro/version/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/pro/version/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/security
/usr/lib/python3/dist-packages/uaclient/api/u/security/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/security/package_manifest
/usr/lib/python3/dist-packages/uaclient/api/u/security/package_manifest/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/security/package_manifest/v1.py
/usr/lib/python3/dist-packages/uaclient/api/u/unattended_upgrades
/usr/lib/python3/dist-packages/uaclient/api/u/unattended_upgrades/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/unattended_upgrades/status
/usr/lib/python3/dist-packages/uaclient/api/u/unattended_upgrades/status/__init__.py
/usr/lib/python3/dist-packages/uaclient/api/u/unattended_upgrades/status/v1.py
/usr/lib/python3/dist-packages/uaclient/apt.py
/usr/lib/python3/dist-packages/uaclient/apt_news.py
/usr/lib/python3/dist-packages/uaclient/cli
/usr/lib/python3/dist-packages/uaclient/cli/__init__.py
/usr/lib/python3/dist-packages/uaclient/cli/api.py
/usr/lib/python3/dist-packages/uaclient/cli/attach.py
/usr/lib/python3/dist-packages/uaclient/cli/auto_attach.py
/usr/lib/python3/dist-packages/uaclient/cli/cli_util.py
/usr/lib/python3/dist-packages/uaclient/cli/collect_logs.py
/usr/lib/python3/dist-packages/uaclient/cli/commands.py
/usr/lib/python3/dist-packages/uaclient/cli/config.py
/usr/lib/python3/dist-packages/uaclient/cli/detach.py
/usr/lib/python3/dist-packages/uaclient/cli/disable.py
/usr/lib/python3/dist-packages/uaclient/cli/enable.py
/usr/lib/python3/dist-packages/uaclient/cli/fix.py
/usr/lib/python3/dist-packages/uaclient/cli/help.py
/usr/lib/python3/dist-packages/uaclient/cli/parser.py
/usr/lib/python3/dist-packages/uaclient/cli/refresh.py
/usr/lib/python3/dist-packages/uaclient/cli/security_status.py
/usr/lib/python3/dist-packages/uaclient/cli/status.py
/usr/lib/python3/dist-packages/uaclient/cli/system.py
/usr/lib/python3/dist-packages/uaclient/clouds
/usr/lib/python3/dist-packages/uaclient/clouds/__init__.py
/usr/lib/python3/dist-packages/uaclient/clouds/aws.py
/usr/lib/python3/dist-packages/uaclient/clouds/azure.py
/usr/lib/python3/dist-packages/uaclient/clouds/gcp.py
/usr/lib/python3/dist-packages/uaclient/clouds/identity.py
/usr/lib/python3/dist-packages/uaclient/config.py
/usr/lib/python3/dist-packages/uaclient/contract.py
/usr/lib/python3/dist-packages/uaclient/contract_data_types.py
/usr/lib/python3/dist-packages/uaclient/daemon
/usr/lib/python3/dist-packages/uaclient/daemon/__init__.py
/usr/lib/python3/dist-packages/uaclient/daemon/poll_for_pro_license.py
/usr/lib/python3/dist-packages/uaclient/daemon/retry_auto_attach.py
/usr/lib/python3/dist-packages/uaclient/data_types.py
/usr/lib/python3/dist-packages/uaclient/defaults.py
/usr/lib/python3/dist-packages/uaclient/entitlements
/usr/lib/python3/dist-packages/uaclient/entitlements/__init__.py
/usr/lib/python3/dist-packages/uaclient/entitlements/anbox.py
/usr/lib/python3/dist-packages/uaclient/entitlements/base.py
/usr/lib/python3/dist-packages/uaclient/entitlements/cc.py
/usr/lib/python3/dist-packages/uaclient/entitlements/cis.py
/usr/lib/python3/dist-packages/uaclient/entitlements/entitlement_status.py
/usr/lib/python3/dist-packages/uaclient/entitlements/esm.py
/usr/lib/python3/dist-packages/uaclient/entitlements/fips.py
/usr/lib/python3/dist-packages/uaclient/entitlements/landscape.py
/usr/lib/python3/dist-packages/uaclient/entitlements/livepatch.py
/usr/lib/python3/dist-packages/uaclient/entitlements/realtime.py
/usr/lib/python3/dist-packages/uaclient/entitlements/repo.py
/usr/lib/python3/dist-packages/uaclient/entitlements/ros.py
/usr/lib/python3/dist-packages/uaclient/event_logger.py
/usr/lib/python3/dist-packages/uaclient/exceptions.py
/usr/lib/python3/dist-packages/uaclient/files
/usr/lib/python3/dist-packages/uaclient/files/__init__.py
/usr/lib/python3/dist-packages/uaclient/files/data_types.py
/usr/lib/python3/dist-packages/uaclient/files/files.py
/usr/lib/python3/dist-packages/uaclient/files/machine_token.py
/usr/lib/python3/dist-packages/uaclient/files/notices.py
/usr/lib/python3/dist-packages/uaclient/files/state_files.py
/usr/lib/python3/dist-packages/uaclient/files/user_config_file.py
/usr/lib/python3/dist-packages/uaclient/gpg.py
/usr/lib/python3/dist-packages/uaclient/http
/usr/lib/python3/dist-packages/uaclient/http/__init__.py
/usr/lib/python3/dist-packages/uaclient/http/serviceclient.py
/usr/lib/python3/dist-packages/uaclient/livepatch.py
/usr/lib/python3/dist-packages/uaclient/lock.py
/usr/lib/python3/dist-packages/uaclient/log.py
/usr/lib/python3/dist-packages/uaclient/messages
/usr/lib/python3/dist-packages/uaclient/messages/__init__.py
/usr/lib/python3/dist-packages/uaclient/messages/urls.py
/usr/lib/python3/dist-packages/uaclient/secret_manager.py
/usr/lib/python3/dist-packages/uaclient/security_status.py
/usr/lib/python3/dist-packages/uaclient/snap.py
/usr/lib/python3/dist-packages/uaclient/status.py
/usr/lib/python3/dist-packages/uaclient/system.py
/usr/lib/python3/dist-packages/uaclient/timer
/usr/lib/python3/dist-packages/uaclient/timer/__init__.py
/usr/lib/python3/dist-packages/uaclient/timer/metering.py
/usr/lib/python3/dist-packages/uaclient/timer/update_contract_info.py
/usr/lib/python3/dist-packages/uaclient/timer/update_messaging.py
/usr/lib/python3/dist-packages/uaclient/types.py
/usr/lib/python3/dist-packages/uaclient/upgrade_lts_contract.py
/usr/lib/python3/dist-packages/uaclient/util.py
/usr/lib/python3/dist-packages/uaclient/version.py
/usr/lib/python3/dist-packages/uaclient/yaml.py
/usr/lib/python3/dist-packages/ubuntu_pro_client-8001.egg-info
/usr/lib/python3/dist-packages/ubuntu_pro_client-8001.egg-info/PKG-INFO
/usr/lib/python3/dist-packages/ubuntu_pro_client-8001.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/ubuntu_pro_client-8001.egg-info/entry_points.txt
/usr/lib/python3/dist-packages/ubuntu_pro_client-8001.egg-info/requires.txt
/usr/lib/python3/dist-packages/ubuntu_pro_client-8001.egg-info/top_level.txt
/usr/lib/ubuntu-advantage
/usr/lib/ubuntu-advantage/apt-esm-json-hook
/usr/lib/ubuntu-advantage/apt_news.py
/usr/lib/ubuntu-advantage/auto_attach.py
/usr/lib/ubuntu-advantage/cloud-id-shim.sh
/usr/lib/ubuntu-advantage/convert_list_to_deb822.py
/usr/lib/ubuntu-advantage/daemon.py
/usr/lib/ubuntu-advantage/esm_cache.py
/usr/lib/ubuntu-advantage/postinst-migrations.sh
/usr/lib/ubuntu-advantage/reboot_cmds.py
/usr/lib/ubuntu-advantage/timer.py
/usr/lib/ubuntu-advantage/upgrade_lts_contract.py
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_ubuntu-advantage-tools.py
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/pro
/usr/share/bash-completion/completions/ua
/usr/share/doc
/usr/share/doc/ubuntu-pro-client
/usr/share/doc/ubuntu-pro-client/changelog.gz
/usr/share/doc/ubuntu-pro-client/copyright
/usr/share/keyrings
/usr/share/keyrings/ubuntu-pro-anbox-cloud.gpg
/usr/share/keyrings/ubuntu-pro-cc-eal.gpg
/usr/share/keyrings/ubuntu-pro-cis.gpg
/usr/share/keyrings/ubuntu-pro-esm-apps.gpg
/usr/share/keyrings/ubuntu-pro-esm-infra.gpg
/usr/share/keyrings/ubuntu-pro-fips-preview.gpg
/usr/share/keyrings/ubuntu-pro-fips.gpg
/usr/share/keyrings/ubuntu-pro-realtime-kernel.gpg
/usr/share/keyrings/ubuntu-pro-ros.gpg
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ubuntu-pro-client
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ubuntu-advantage.1.gz
/var
/var/lib
/var/lib/ubuntu-advantage
/usr/bin/pro
/usr/bin/ua
/usr/share/man/man1/pro.1.gz
/usr/share/man/man1/ua.1.gz
