# dietpi/var/lib/dpkg/info/cron.conffiles
/etc/default/cron
/etc/init.d/cron
/etc/pam.d/cron
