# dietpi/var/lib/dpkg/info/cron-daemon-common.list
/.
/etc
/etc/cron.d
/etc/cron.d/.placeholder
/etc/cron.daily
/etc/cron.daily/.placeholder
/etc/cron.hourly
/etc/cron.hourly/.placeholder
/etc/cron.monthly
/etc/cron.monthly/.placeholder
/etc/cron.weekly
/etc/cron.weekly/.placeholder
/etc/cron.yearly
/etc/cron.yearly/.placeholder
/etc/crontab
/usr
/usr/share
/usr/share/doc
/usr/share/doc/cron-daemon-common
/usr/share/doc/cron-daemon-common/changelog.Debian.gz
/usr/share/doc/cron-daemon-common/changelog.gz
/usr/share/doc/cron-daemon-common/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/cron-daemon-common
/var
/var/spool
/var/spool/cron
/var/spool/cron/crontabs
