# dietpi/var/lib/dpkg/info/dropbear.list
/.
/etc
/etc/default
/etc/default/dropbear
/etc/dropbear
/etc/dropbear/log
/etc/dropbear/log/run
/etc/dropbear/run
/etc/init.d
/etc/init.d/dropbear
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dropbear.service
/usr
/usr/share
/usr/share/doc
/usr/share/doc/dropbear
/usr/share/doc/dropbear/NEWS.Debian.gz
/usr/share/doc/dropbear/README.Debian
/usr/share/doc/dropbear/README.runit
/usr/share/doc/dropbear/changelog.Debian.gz
/usr/share/doc/dropbear/changelog.gz
/usr/share/doc/dropbear/copyright
/etc/dropbear/log/main
