# dietpi/var/lib/dpkg/info/dropbear.conffiles
/etc/default/dropbear
/etc/dropbear/log/run
/etc/dropbear/run
/etc/init.d/dropbear
