# dietpi/var/lib/dpkg/info/ifupdown.conffiles
/etc/default/networking
/etc/init.d/networking
/etc/network/if-down.d/resolved
/etc/network/if-up.d/resolved
